# systemd/wg-captive-agent.service
[Unit]
Description=WireGuard captive firewall agent
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
EnvironmentFile=/etc/wg-captive-agent.env
ExecStart=/usr/local/sbin/wg-captive-agent run
ExecStopPost=/usr/local/sbin/wg-captive-agent cleanup
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

